# templates/gunicorn.socket.j2
[Unit]
Description=gunicorn socket for {{ gunicorn-django-project-name }}

[Socket]
ListenStream=/tmp/gunicorn.{{ gunicorn-django-project-name }}.sock

[Install]
WantedBy=sockets.target
