# templates/gunicorn.socket.j2
[Unit]
Description=gunicorn socket for {{ gunicorn_django_project_name }}

[Socket]
ListenStream=/tmp/gunicorn.{{ gunicorn_django_project_name }}.sock

[Install]
WantedBy=sockets.target
