# templates/gunicorn.socket.j2
[Unit]
Description=gunicorn socket for {{ gunicorn_django_project_name }}.%I

[Socket]
ListenStream=/var/run/gunicorn.{{ gunicorn_django_project_name }}.%i.sock

[Install]
WantedBy=sockets.target
